# docker_root/etc/services.d/unifi-bouncer/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    /app/crowdsec-custom-bouncer -c /etc/crowdsec/bouncers/crowdsec-custom-bouncer.yaml
